# deploy/nik-update.service
[Unit]
Description=Nik self-update check

[Service]
Type=oneshot
User=nik
ExecStart=%h/git/deploy/update.sh
Environment=PATH=/usr/local/go/bin:/usr/bin:/bin
Environment=HOME=%h
