# deploy/nik.service
[Unit]
Description=Nik - Noetic Intelligence Kernel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nik
ExecStart=%h/git/nik --home %h/workspace
WorkingDirectory=%h/git
Restart=always
RestartSec=5
Environment=CGO_ENABLED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
