# internal/daemonctl/nikd.service.tmpl
[Unit]
Description=nikd
After=network-online.target
Wants=network-online.target
# Crash-loop guard: if the daemon dies 10 times in 2 minutes systemd gives
# up. Recover with `systemctl --user reset-failed nikd.service && systemctl --user start nikd.service`.
StartLimitBurst=10
StartLimitIntervalSec=120

[Service]
Type=simple
ExecStart={{ .NikBinary }} daemon --home {{ .NikHome }}
# Always restart: the TUI signals SIGTERM to the daemon to apply config
# changes, which is a clean exit (status=0). Restart=on-failure would skip
# that case. See setup.go updateDone.
Restart=always
RestartSec=5
Environment=NIK_HOME={{ .NikHome }}

[Install]
WantedBy=default.target
